# ausfallsicherheit/ansible/roles/minikube/templates/minikube.service.j2
[Unit]
Description=minikube

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStart=minikube start --driver=docker --force
ExecStop=minikube stop

[Install]
WantedBy=multi-user.target
